# el9/rpmbuild/BUILD/jobarranger-6.1.1/misc/init.d/redhat/jobarg-agentd
#!/bin/bash
#
# jobarg-agentd Startup script for the Job Arranger Agentd
#
# chkconfig: - 90 15
# description:  Job Arranger Agentd
# processname: jobarg_agentd
# config: /etc/jobarranger/jobarg_agentd.conf
# pidfile: /var/run/jobarranger/jobarg_agentd.pid

# jobarg details
JOBARG=/usr/sbin/jobarg_agentd
CONF=/etc/jobarranger/jobarg_agentd.conf
PIDFILE=/var/run/jobarranger/jobarg_agentd.pid

# Source function library.
. /etc/rc.d/init.d/functions

# Source networking configuration.
. /etc/sysconfig/network

# Check that networking is up.
[ ${NETWORKING} = "no" ] && exit 0

[ -e $CONF ] || exit 6

RETVAL=0

case "$1" in
    start)
        echo -n "Starting Job Arranger Agentd: "
        daemon $JOBARG -c $CONF
        RETVAL=$?
        echo
        [ $RETVAL -eq 0 ] && touch /var/lock/subsys/jobarg_agentd
        ;;
    stop)
        echo -n "Shutting down Job Arranger Agentd: "
        killproc -p $PIDFILE
        RETVAL=$?
        echo
        [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/jobarg_agentd
        ;;
    restart)
        $0 stop
        sleep 5
        $0 start
        RETVAL=$?
        ;;
    condrestart)
        if [ -f /var/lock/subsys/jobarg_agentd ]; then
            $0 stop
            $0 start
        fi
        RETVAL=$?
        ;;
    status)
        status $JOBARG
        RETVAL=$?
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|condrestart|status}"
        exit 1
        ;;
esac

exit $RETVAL
