# el9/rpmbuild/SOURCES/jobarranger-6.1.1/.svn/pristine/71/7176964b382b880a04e8233685701089075dd863.svn-base
#!/bin/bash
#
# jobarg-server and jobarg-monitor Startup script for the Job Arranger Server
#
# chkconfig: - 90 15
# description:  Job Arranger Server and Job Arranger Monitor
# processname: jobarg_server, jobarg_monitor

# Source networking configuration.
. /etc/sysconfig/network

# Check that networking is up.
[ ${NETWORKING} = "no" ] && exit 0

RETVAL=0

case "$1" in
    start)
        /etc/rc.d/init.d/jobarg-server $1
        RETVAL=$?
        if [ $RETVAL -ne 0 ]; then
            exit $RETVAL
        fi
        /etc/rc.d/init.d/jobarg-monitor $1
        RETVAL=$?
        ;;
    stop)
        /etc/rc.d/init.d/jobarg-monitor $1
        /etc/rc.d/init.d/jobarg-server $1
        RETVAL=$?
        ;;
    restart)
        $0 stop
        sleep 3
        $0 start
        RETVAL=$?
        ;;
    *)
        echo "Usage: $0 {start|stop|restart}"
        exit 1
        ;;
esac

exit $RETVAL
